# conf.d/fdm.conf
[program:fdm]
command=/opt/freedownloadmanager/fdm -platform vnc
autorestart=true